# etc/crontab
0,15,30,45 *     * * *     packages-www cd /var/www/packages.grml.org && systemd-cat --identifier=packages-www-cron flock /var/www/packages.grml.org/git/lock ./packages.grml.org/update.rb
